# com.user.apple-reminders-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.apple-reminders-bridge</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/basu/CodeBase/apps/set_reminders/apple_reminders_bridge.py</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <false/>
    
    <key>StandardOutPath</key>
    <string>/Users/basu/Library/Logs/apple-reminders-bridge.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/basu/Library/Logs/apple-reminders-bridge-error.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/basu/CodeBase/apps/set_reminders</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
